# supervisord.conf
[supervisord]
nodaemon=true

[program:llama2_server]
command=uvicorn run_llama2_server:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
stderr_logfile=/var/log/llama2.err.log
stdout_logfile=/var/log/llama2.out.log

[program:streamlit_app]
command=streamlit run chatbot.py --server.port 8501
autostart=true
autorestart=true
stderr_logfile=/var/log/streamlit.err.log
stdout_logfile=/var/log/streamlit.out.log
